# auto_commit.service
[Unit]
Description=Commandline Utility that automates git commits

[Service]
ExecStart=/home/codam/commit/commit.sh start

[Install]
WantedBy=multi-user.target
